# layout/DEBIAN/postrm
#!/bin/sh

echo Reloading daemon
if launchctl unload -wF /System/Library/LaunchDaemons/com.apple.TextInput.kbd.plist 2>/dev/null; then
  launchctl load -wF /System/Library/LaunchDaemons/com.apple.TextInput.kbd.plist 2>/dev/null || true
else
  # Elucubratus launchctl may fail, fallback to killall
  killall -TERM kbd 2>/dev/null || true
fi

exit 0
